# dockerfile/spoter-supervisor.conf
[program:spoter-controller]
command=/home/spoter/run-spoter-controller.sh
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stopasgroup=true
stopsignal=QUIT
